# deploy/systemd/dwarvenpick.service
[Unit]
Description=dwarvenpick OpenSQL Workbench Backend
After=network.target

[Service]
Type=simple
User=dwarvenpick
Group=dwarvenpick
WorkingDirectory=/opt/dwarvenpick
Environment=SPRING_PROFILES_ACTIVE=prod
EnvironmentFile=-/etc/dwarvenpick/dwarvenpick.env
ExecStart=/usr/bin/java -jar /opt/dwarvenpick/app.jar
Restart=on-failure
RestartSec=5
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
